# activitywatch-exporter.timer
[Unit]
Description=run activitywatch-exporter daily

[Timer]
OnCalendar=Daily
Persistent=true

[Install]
WantedBy=timers.target